# system/ivy-boot.service
[Unit]
Description=Ivy boot script

[Service]
Type=oneshot
User=root
StandardOutput=journal+console
StandardError=journal+console
StandardInput=tty
ExecStop=/usr/bin/python3.6 /home/ivy/ivy/system/on-boot.py
TimeoutSec=infinity

[Install]
WantedBy=network-online.target
